# templates/default/unicorn_rails.init.erb
#! /bin/bash

# Generated by Chef for <%= @fqdn %>
# Local modifications will be overwritten.

# unicorn_rails init script for Linux 
#
# chkconfig: 2345 60 40
# description:  unicorn_rails init script for Gitlab
# processname: unicorn_rails

### BEGIN INIT INFO
# Provides:          unicorn_rails
# Required-Start:    $local_fs $remote_fs $network $syslog
# Required-Stop:     $local_fs $remote_fs $network $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: unicorn_rails init script
# Description:       unicorn_rails init script as shipped with gitlab cookbook
### END INIT INFO

DAEMON_OPTS="-c <%= @gitlab_app_home %>/config/unicorn.rb -E production -D"
NAME=unicorn_rails
DESC=gitlab
PID=<%= @gitlab_app_home %>/tmp/pids/unicorn.pid
RESQUE_PID=<%= @gitlab_app_home %>/tmp/pids/resque_worker.pid

case "$1" in
  start)
        if [ -x /etc/profile.d/gitlab.sh ]; then
          . /etc/profile.d/gitlab.sh
        fi
        CD_TO_APP_DIR="cd <%= @gitlab_app_home %>"
        START_DAEMON_PROCESS="bundle exec unicorn_rails $DAEMON_OPTS"
        START_RESQUE_PROCESS="./resque.sh"

        echo -n "Starting $DESC: "
        if [ `whoami` = root ]; then
          sudo -u gitlab sh -l -c "$CD_TO_APP_DIR > /dev/null 2>&1 && $START_DAEMON_PROCESS && $START_RESQUE_PROCESS"
        else
          $CD_TO_APP_DIR > /dev/null 2>&1 && $START_DAEMON_PROCESS && $START_RESQUE_PROCESS
        fi
        echo "$NAME."
        ;;
  stop)
        echo -n "Stopping $DESC: "
        kill -QUIT `cat $PID`
        kill -QUIT `cat $RESQUE_PID`
        echo "$NAME."
        ;;
  restart)
        echo -n "Restarting $DESC: "
        kill -USR2 `cat $PID`
        kill -USR2 `cat $RESQUE_PID`
        echo "$NAME."
        ;;
  reload)
        echo -n "Reloading $DESC configuration: "
        kill -HUP `cat $PID`
        kill -HUP `cat $RESQUE_PID`
        echo "$NAME."
        ;;
  *)
        echo "Usage: $NAME {start|stop|restart|reload}" >&2
        exit 1
        ;;
esac

exit 0
